# docker/php/cron/crontab
* */3 * * * php /var/www/html/cron.php >> /dev/null 2>&1


